# etc/supervisord.conf
[program:david]
command=/srv/usr/david/app/venv/bin/gunicorn app.py:app -c gunicorn.conf.py
directory=/srv/user/david/app/
user=nobody
autostart=true
autorestart=true
redirect_stderr=True
